# stage1/05-firmware/files/iti-boot-config.timer
[Unit]
Description=Run iti boot config on every boot

[Timer]
OnBootSec=5sec
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target