# stage1/05-firmware/files/iti-boot-config.service
[Unit]
Description=Iti Boot Config

[Service]
Type=oneshot
ExecStart=/usr/local/bin/iti-boot-config
RemainAfterExit=true