# stage1/05-firmware/files/iti-boot-config.timer
[Unit]
Description=Run iti boot config on every boot

[Timer]
OnBootSec=15sec
OnUnitActiveSec=10sec

[Install]
WantedBy=timers.target