# stage1/07-iti-good-indicator/files/iti-lcd-good.service
[Unit]
Description=Display hardware-good indicator on LCD

[Service]
Type=oneshot
ExecStart=/usr/local/bin/iti-lcd-good
RemainAfterExit=true
