# stage1/07-iti-good-indicator/files/iti-lcd-good.timer
[Unit]
Description=Check if the battery is well seated

[Timer]
OnBootSec=5sec
OnUnitActiveSec=5sec

[Install]
WantedBy=timers.target
